# roles/python_app/templates/log-stats.service.j2
[Unit]
Description=System monitor metrics sampler (one-shot)

[Service]
Type=oneshot
EnvironmentFile=/opt/system-monitor/.env
ExecStart=/opt/system-monitor/venv/bin/python /opt/system-monitor/log_stats.py
