# roles/python_app/templates/log-stats.timer.j2
[Unit]
Description=Run log-stats every 15 seconds

[Timer]
OnUnitActiveSec=15s
AccuracySec=1s
Unit=log-stats.service

[Install]
WantedBy=timers.target
