# roles/python_app/templates/system-monitor.service.j2
[Unit]
Description=System Monitor (writes metrics to DB)
After=network.target

[Service]
EnvironmentFile=/opt/system-monitor/.env
ExecStart=/opt/system-monitor/venv/bin/python /opt/system-monitor/log_stats.py
Type=oneshot
